# scripts/com.mls.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd agent: rebuild the MLS webapp dashboard daily at 07:00 local.
  Install:
    cp scripts/com.mls.dashboard.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.mls.dashboard.plist
  Run once now:
    launchctl start com.mls.dashboard
  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.mls.dashboard.plist
  Logs land in logs/daily_build.{out,err}.log under the repo.
  Adjust the repo path below if your checkout differs.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mls.dashboard</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/ryangerda/Development/MLS/scripts/daily_build.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>7</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/ryangerda/Development/MLS</string>

  <key>StandardOutPath</key>
  <string>/Users/ryangerda/Development/MLS/logs/daily_build.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/ryangerda/Development/MLS/logs/daily_build.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
